# roles/pycharm/templates/etc/profile.d/pycharm.sh.j2
PYCHARM_DIR="{{pycharm_install_path}}/pycharm-{{pycharm_type}}-{{pycharm_version}}"
PATH="${PYCHARM_DIR}/bin:${PATH}"
